# deploy/celerybeat.service
[Unit]
Description=StatementHub Celery Beat Scheduler
After=network.target redis.service postgresql.service
Requires=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/statementhub
Environment="PATH=/opt/statementhub/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/statementhub/.env
ExecStart=/opt/statementhub/venv/bin/celery -A config beat \
    --loglevel=info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler \
    --pidfile=/var/run/celerybeat.pid
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/statementhub/celery-beat.log
StandardError=append:/var/log/statementhub/celery-beat.log
# Hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
